# crontab
# 手动模式下，每5分钟把缓存存进数据库
*/5 * * * * curl -d "" "http://127.0.0.1:8080/analyze/finish"
# 心跳
59 23 * * * curl "http://127.0.0.1:8080/analyze/heartbeat"
# Test
# */1 * * * * echo "hello" >> /tmp/log
